# scripts/com.agenttoolbox.broker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Keep the agent-toolbox broker running at login. Install:
    cp scripts/com.agenttoolbox.broker.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.agenttoolbox.broker.plist
  Logs: /tmp/agtbx-broker.log   Stop: launchctl unload ~/Library/LaunchAgents/com.agenttoolbox.broker.plist

  Runs the broker through a login shell (`bash -lc`) so the user's profile sets
  PATH (node/nvm, Android SDK) and $HOME expands at runtime — no machine-specific
  paths baked in. Assumes the repo lives at $HOME/code/agent-toolbox; edit the
  command if yours is elsewhere.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agenttoolbox.broker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec "$HOME/code/agent-toolbox/scripts/broker.sh"</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/agtbx-broker.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/agtbx-broker.log</string>
</dict>
</plist>
